# deploy/scripts/systemd/ronin-backup.service
[Unit]
Description=Ronin backup (SQLite databases, replays)

[Service]
Type=oneshot
ExecStart=/opt/ronin/scripts/backup.sh
